# apps/api-gateway/edge/zhilian-edge-shokz.service
[Unit]
Description=Zhilian OS Raspberry Pi Shokz Callback Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zhilian-edge
Group=zhilian-edge
EnvironmentFile=/etc/zhilian-edge/edge-node.env
WorkingDirectory=/opt/zhilian-edge
ExecStart=/usr/bin/python3 /opt/zhilian-edge/shokz_callback_daemon.py
Restart=always
RestartSec=5
StateDirectory=zhilian-edge
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
